# config/syncronix-marketing.service
[Unit]
Description=Syncronix Marketing Automation (Hotmart -> WhatsApp)
After=network.target

[Service]
User=whatsapp-agent
Group=whatsapp-agent
WorkingDirectory=/root/.openclaw/workspace/whatsapp-ai-agent
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/root/.openclaw/workspace/whatsapp-ai-agent/.env
Environment="MARKETING_AUTOMATION_ENABLED=true"
Environment="AI_AGENT_ENABLED=false"
ExecStart=/usr/bin/python3 -m uvicorn src.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=5
KillMode=mixed
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
